# deploy/bestweather.service
# systemd unit for AltijdMooiWeer on the altijdmooiweer.nl server.
# Install to /etc/systemd/system/bestweather.service.
#
# Isolated by design: binds ONLY to 127.0.0.1:8800; nginx reverse-proxies it for
# the altijdmooiweer.nl vhost. Runs as the deploy user, matching the house style of
# gunicorn-ehbo.service on this host.

[Unit]
Description=AltijdMooiWeer (altijdmooiweer.nl) ASGI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/altijdmooiweer.nl
EnvironmentFile=-/var/www/altijdmooiweer.nl/.env
ExecStart=/var/www/altijdmooiweer.nl/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8800 --workers 1
Restart=on-failure
RestartSec=3

NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
